# configs/docker/supervisor-psono-worker.conf
[program:crond]
command = /usr/sbin/cron -f
autostart = true