# configs/docker/supervisor-psono-server.conf
[program:nginx]
command = /usr/sbin/nginx

[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /root/configs/nginx/psono_uwsgi.ini
