# configs/docker/supervisor-psono-worker.conf
[program:crond]
command = /usr/sbin/cron
autostart = true